# deploy/systemd/container-cleanup-daily.service
[Unit]
Description=Daily Docker builder cleanup with Telegram report
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/container-cleanup/container_cleanup.sh daily 85
