# deploy/systemd/container-cleanup-weekly.timer
[Unit]
Description=Run weekly Docker image and volume cleanup when disk usage is high

[Timer]
OnCalendar=Sun 04:40:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
