# deploy/systemd/container-cleanup-daily.timer
[Unit]
Description=Run daily Docker builder cleanup

[Timer]
OnCalendar=*-*-* 04:20:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target

